# config/docker/php/crontab
* * * * * www-data /usr/bin/php bin/console app:cleanup
